# deploy/rumors.service
[Unit]
Description=Field Intel — rumors API (uvicorn via gunicorn)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/rumors/backend
EnvironmentFile=/opt/rumors/.env
ExecStart=/opt/rumors/venv/bin/gunicorn app.main:app \
    --preload \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 2 \
    --bind 127.0.0.1:5003 \
    --access-logfile /var/log/rumors-access.log \
    --error-logfile /var/log/rumors-error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
